# exoworker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.exo.worker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Exo.app/Contents/MacOS/exo</string>
        <string>--node-name</string>
        <string>Node #</string>
        <string>--node-host</string>
        <string>10.0.0.#</string>
        <string>--memory-limit-gb</string>
        <string>6</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>EXO_LIBP2P_NAMESPACE</key>
        <string>soit_exo_cluster</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/exo.node#.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/exo.node#.err.log</string>
</dict>
</plist>
